# t-rex-miner-img/t-rex-miner.service
[Unit]
Description=T-Rex Miner Service

[Service]
ExecStart=/opt/start-mining.sh

[Install]
WantedBy=multi-user.target